# nginx/preconfig/etc/s6-overlay/s6-rc.d/init-folders/run
#!/command/execlineb -P
mkdir -p \
    /config/ssl \
    /config/nginx
