# nginx/preconfig/etc/s6-overlay/s6-rc.d/init-folders/run
#!/bin/ash
mkdir -p \
    /config/ssl \
    /config/ssl/default \
    /config/nginx \
    /config/nginx/sites \
    /config/www
