# nginx/preconfig/etc/s6-overlay/s6-rc.d/init-files/run
#!/bin/ash
if [ ! -f /config/nginx/nginx.conf ]; then
    cp /defaults/nginx.conf /config/nginx/nginx.conf
fi